# systemd/ukabu-unjail.timer
[Unit]
Description=UKABU WAF Unjail Timer (Every 5 Minutes)
[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
[Install]
WantedBy=timers.target
