# systemd/ukabu-unjail.service
[Unit]
Description=UKABU WAF Unjail Expired Blocks
[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/ukabu-unjail.sh
